# client/ebox/debian/ebox.cron.hourly
#!/bin/bash

/etc/init.d/ebox events restart
