# client/printers/debian/ebox-printers.postinst
#!/bin/bash

#DEBHELPER#

onInstall() {
    perl -e 'use EBox::Printers; EBox::Printers::onInstall();'
} 


case "$1" in
    configure)
        kill `pidof gconfd-2` >/dev/null 2>&1 || true

        onInstall

        # Restart eBox apache
        invoke-rc.d ebox apache restart || true

     ;;
esac



exit 0
